# rootfs/etc/supervisor/conf.d/auto-create-user.conf
[program:auto-create-user]
command=/usr/bin/auto-create-user.py
autostart=false
user=root
autorestart=true
priority=100
